# project/Crawl_Linux_Cmd/Download/系统管理/crontab
<!DOCTYPE html>
<html>
<head>
<meta charset="UTF-8" />
<meta name="renderer" content="webkit">
<title>
crontab命令_Linux crontab 命令用法详解：提交和管理用户的需要周期性执行的任务</title>
<link rel="shortcut icon" href="http://man.linuxde.net/favicon.ico" type="image/x-icon" />
<link rel="pingback" href="http://man.linuxde.net/xmlrpc.php" />
<link href="http://112.124.15.130/res/linuxde_man/style.css" rel="stylesheet" type="text/css">
<script type="text/javascript" src="http://112.124.15.130/res/linuxde_man/js/jquery-1.8.3.min.js"></script>
<!--[if IE 6]>
<script src="http://112.124.15.130/res/linuxde_man/js/belatedPNG.js"></script>
<script>DD_belatedPNG.fix('*');</script>
<script src="http://letskillie6.googlecode.com/svn/trunk/letskillie6.zh_CN.pack.js"></script>
<![endif]-->
<link rel='canonical' href='http://man.linuxde.net/crontab' />
</head>
<body>
<noscript>
<p class="nojs">您的浏览器已经禁用了脚本，这会严重影响您正常使用本站的功能，请开启！</p>
</noscript>
<div class="man_linuxde_net_page">
<div class="layout_hd">
    <div class="logo l">
        <h2><a title="Linux命令大全" href="http://man.linuxde.net/">
            Linux命令大全            </a></h2>
    </div>
    <div class="nav r"><ul id="menu-%e8%8f%9c%e5%8d%951" class="menu"><li id="menu-item-4" class="menu-item menu-item-type-custom menu-item-object-custom menu-item-home menu-item-4"><a href="http://man.linuxde.net/">首页搜索</a></li>
<li id="menu-item-6" class="menu-item menu-item-type-taxonomy menu-item-object-category current-post-ancestor current-menu-parent current-post-parent menu-item-6"><a href="http://man.linuxde.net/par/2">系统管理</a></li>
<li id="menu-item-7" class="menu-item menu-item-type-taxonomy menu-item-object-category menu-item-7"><a href="http://man.linuxde.net/par/5">网络管理</a></li>
<li id="menu-item-8" class="menu-item menu-item-type-taxonomy menu-item-object-category menu-item-8"><a href="http://man.linuxde.net/par/4">软件·打印·开发·工具</a></li>
<li id="menu-item-9" class="menu-item menu-item-type-taxonomy menu-item-object-category menu-item-9"><a href="http://man.linuxde.net/par/1">文件和目录管理</a></li>
<li id="menu-item-10" class="menu-item menu-item-type-taxonomy menu-item-object-category menu-item-10"><a href="http://man.linuxde.net/par/3">硬件·内核·Shell·监测</a></li>
<li id="menu-item-2096" class="menu-item menu-item-type-custom menu-item-object-custom menu-item-2096"><a target="_blank" href="http://www.linuxde.net">Linux Today</a></li>
</ul></div>    <div class="c"></div>
</div>
<div class="layout layout_bd">
<script type='text/javascript'>
function addLink() {
    var body_element = document.getElementsByTagName('body')[0];
    var selection;
    selection = window.getSelection();
    var pagelink = "<br /><br /> 来自: <a href='"+document.location.href+"'>"+document.location.href+"</a>"; 
    var copy_text = selection + pagelink;
    var new_div = document.createElement('div');
    new_div.style.left='-99999px';
    new_div.style.position='absolute';
    body_element.appendChild(new_div );
    new_div.innerHTML = copy_text ;
    selection.selectAllChildren(new_div );
    window.setTimeout(function() {
        body_element.removeChild(new_div );
    },0);
}
document.oncopy = addLink;
</script>
<div class="single_default">
        <div class="main">
        <div class="breadcrumb">当前位置：<a href="http://man.linuxde.net">首页</a> &raquo; <a href="http://man.linuxde.net/par/2" rel="category tag">系统管理</a> &raquo; crontab</div>
        <div class="post_hd clearfix">
            <h1 class="l">crontab命令</h1>
            <div class="tag f_yh r">
                <a href="http://man.linuxde.net/sub/%e8%bf%9b%e7%a8%8b%e5%92%8c%e4%bd%9c%e4%b8%9a%e7%ae%a1%e7%90%86" rel="tag">进程和作业管理</a>            </div>
        </div>
        <div class="post_bd post">
            <p>
	<strong>crontab命令</strong>被用来提交和管理用户的需要周期性执行的任务，与windows下的计划任务类似，当安装完成操作系统后，默认会安装此服务工具，并且会自动启动crond进程，crond进程每分钟会定期检查是否有要执行的任务，如果有要执行的任务，则自动执行该任务。
</p>
<h3>
	语法<br />
</h3>
<pre>
crontab(选项)(参数)</pre>
<h3>
	选项<br />
</h3>
<pre>
-e：编辑该用户的计时器设置；
-l：列出该用户的计时器设置；
-r：删除该用户的计时器设置；
-u&lt;用户名称&gt;：指定要设定计时器的用户名称。</pre>
<h3>
	参数<br />
</h3>
<p>
	crontab文件：指定包含待执行任务的crontab文件。
</p>
<h3>
	知识扩展<br />
</h3>
<p>
	Linux下的任务调度分为两类：<strong>系统任务调度</strong>和<strong>用户任务调度</strong>。
</p>
<p>
	<strong>系统任务调度：</strong>系统周期性所要执行的工作，比如写缓存数据到硬盘、日志清理等。在<code>/etc</code>目录下有一个crontab文件，这个就是系统任务调度的配置文件。
</p>
<p>
	<code>/etc/crontab</code>文件包括下面几行：
</p>
<pre>
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=&quot;&quot;HOME=/

# run-parts
51 * * * * root run-parts /etc/cron.hourly
24 7 * * * root run-parts /etc/cron.daily
22 4 * * 0 root run-parts /etc/cron.weekly
42 4 1 * * root run-parts /etc/cron.monthly</pre>
<p>
	前四行是用来配置crond任务运行的环境变量，第一行SHELL变量指定了系统要使用哪个shell，这里是bash，第二行PATH变量指定了系统执行命令的路径，第三行MAILTO变量指定了crond的任务执行信息将通过电子邮件发送给root用户，如果MAILTO变量的值为空，则表示不发送任务执行信息给用户，第四行的HOME变量指定了在执行命令或者脚本时使用的主目录。
</p>
<p>
	<strong>用户任务调度：</strong>用户定期要执行的工作，比如用户数据备份、定时邮件提醒等。用户可以使用&nbsp;crontab&nbsp;工具来定制自己的计划任务。所有用户定义的crontab文件都被保存在<code>/var/spool/cron</code>目录中。其文件名与用户名一致，使用者权限文件如下：
</p>
<pre>
/etc/cron.deny     该文件中所列用户不允许使用crontab命令
/etc/cron.allow    该文件中所列用户允许使用crontab命令
/var/spool/cron/   所有用户crontab文件存放的目录,以用户名命名</pre>
<p>
	crontab文件的含义：用户所建立的crontab文件中，每一行都代表一项任务，每行的每个字段代表一项设置，它的格式共分为六个字段，前五段是时间设定段，第六段是要执行的命令段，格式如下：
</p>
<pre>
minute   hour   day   month   week   <span class='wp_keywordlink'><a href="http://man.linuxde.net/command" title="command命令" target="_blank">command</a></span>     <span style="color:#FF0000;">顺序：分 时 日 月 周</span></pre>
<p>
	其中：
</p>
<ul>
<li>
		minute：&nbsp;表示分钟，可以是从0到59之间的任何整数。
	</li>
<li>
		hour：表示小时，可以是从0到23之间的任何整数。
	</li>
<li>
		day：表示日期，可以是从1到31之间的任何整数。
	</li>
<li>
		month：表示月份，可以是从1到12之间的任何整数。
	</li>
<li>
		week：表示星期几，可以是从0到7之间的任何整数，这里的0或7代表星期日。
	</li>
<li>
		command：要执行的命令，可以是系统命令，也可以是自己编写的脚本文件。
	</li>
</ul>
<p>
	在以上各个字段中，还可以使用以下特殊字符：
</p>
<ul>
<li>
		星号（*）：代表所有可能的值，例如month字段如果是星号，则表示在满足其它字段的制约条件后每月都执行该命令操作。
	</li>
<li>
		逗号（,）：可以用逗号隔开的值指定一个列表范围，例如，&ldquo;1,2,5,7,8,9&rdquo;
	</li>
<li>
		中杠（-）：可以用整数之间的中杠表示一个整数范围，例如&ldquo;2-6&rdquo;表示&ldquo;2,3,4,5,6&rdquo;
	</li>
<li>
		正斜线（/）：可以用正斜线指定时间的间隔频率，例如&ldquo;0-23/2&rdquo;表示每两小时执行一次。同时正斜线可以和星号一起使用，例如*/10，如果用在minute字段，表示每十分钟执行一次。
	</li>
</ul>
<p>
	<strong>crond服务</strong>
</p>
<pre>
/sbin/<span class='wp_keywordlink'><a href="http://man.linuxde.net/service" title="service命令" target="_blank">service</a></span> crond start    //启动服务
/sbin/service crond stop     //关闭服务
/sbin/service crond restart  //重启服务
/sbin/service crond reload   //重新载入配置</pre>
<p>
	查看crontab服务状态：
</p>
<pre>
service crond status</pre>
<p>
	手动启动crontab服务：
</p>
<pre>
service crond start</pre>
<p>
	查看crontab服务是否已设置为开机启动，执行命令：
</p>
<pre>
<span class='wp_keywordlink'><a href="http://man.linuxde.net/ntsysv" title="ntsysv命令" target="_blank">ntsysv</a></span></pre>
<p>
	加入开机自动启动：
</p>
<pre>
<span class='wp_keywordlink'><a href="http://man.linuxde.net/chkconfig" title="chkconfig命令" target="_blank">chkconfig</a></span> &ndash;level 35 crond on</pre>
<h3>
	实例<br />
</h3>
<p>
	每1分钟执行一次command
</p>
<pre>
* * * * * command</pre>
<p>
	每小时的第3和第15分钟执行
</p>
<pre>
3,15 * * * * command</pre>
<p>
	在上午8点到11点的第3和第15分钟执行
</p>
<pre>
3,15 8-11 * * * command</pre>
<p>
	每隔两天的上午8点到11点的第3和第15分钟执行
</p>
<pre>
3,15 8-11 */2 * * command</pre>
<p>
	每个星期一的上午8点到11点的第3和第15分钟执行
</p>
<pre>
3,15 8-11 * * 1 command</pre>
<p>
	每晚的21:30重启smb&nbsp;
</p>
<pre>
30 21 * * * /etc/<span class='wp_keywordlink'><a href="http://man.linuxde.net/init" title="init命令" target="_blank">init</a></span>.d/smb restart</pre>
<p>
	每月1、10、22日的4&nbsp;:&nbsp;45重启smb&nbsp;
</p>
<pre>
45 4 1,10,22 * * /etc/init.d/smb restart</pre>
<p>
	每周六、周日的1:10重启smb
</p>
<pre>
10 1 * * 6,0 /etc/init.d/smb restart</pre>
<p>
	每天18&nbsp;:&nbsp;00至23&nbsp;:&nbsp;00之间每隔30分钟重启smb&nbsp;
</p>
<pre>
0,30 18-23 * * * /etc/init.d/smb restart</pre>
<p>
	每星期六的晚上11:00&nbsp;pm重启smb&nbsp;
</p>
<pre>
0 23 * * 6 /etc/init.d/smb restart</pre>
<p>
	每一小时重启smb&nbsp;
</p>
<pre>
* */1 * * * /etc/init.d/smb restart</pre>
<p>
	晚上11点到早上7点之间，每隔一小时重启smb
</p>
<pre>
* 23-7/1 * * * /etc/init.d/smb restart</pre>
<p>
	每月的4号与每周一到周三的11点重启smb&nbsp;
</p>
<pre>
0 11 4 * mon-wed /etc/init.d/smb restart</pre>
<p>
	一月一号的4点重启smb
</p>
<pre>
0 4 1 jan * /etc/init.d/smb restart</pre>
<p>
	每小时执行<code>/etc/cron.hourly</code>目录内的脚本
</p>
<pre>
01 * * * * root run-parts /etc/cron.hourly</pre>
        </div>
        <div class="post_ft">
            <div class="title">最近更新的命令</div>
            <ul>
                <li><a href="http://man.linuxde.net/cut" title="cut">cut</a></li><li><a href="http://man.linuxde.net/sh" title="sh">sh</a></li><li><a href="http://man.linuxde.net/seq" title="seq">seq</a></li><li><a href="http://man.linuxde.net/awk" title="awk">awk</a></li><li><a href="http://man.linuxde.net/syslog" title="syslog">syslog</a></li><li><a href="http://man.linuxde.net/inotifywait" title="inotifywait">inotifywait</a></li><li><a href="http://man.linuxde.net/dd" title="dd">dd</a></li><li><a href="http://man.linuxde.net/read" title="read">read</a></li><li><a href="http://man.linuxde.net/uname" title="uname">uname</a></li><li><a href="http://man.linuxde.net/pssh" title="pssh">pssh</a></li><li><a href="http://man.linuxde.net/tar" title="tar">tar</a></li><li><a href="http://man.linuxde.net/axel" title="axel">axel</a></li><li><a href="http://man.linuxde.net/losetup" title="losetup">losetup</a></li><li><a href="http://man.linuxde.net/lsb_release" title="lsb_release">lsb_release</a></li><li><a href="http://man.linuxde.net/tcpreplay" title="tcpreplay">tcpreplay</a></li><li><a href="http://man.linuxde.net/strings" title="strings">strings</a></li><li><a href="http://man.linuxde.net/screen" title="screen">screen</a></li><li><a href="http://man.linuxde.net/speedtest-cli" title="speedtest-cli">speedtest-cli</a></li><li><a href="http://man.linuxde.net/clockdiff" title="clockdiff">clockdiff</a></li><li><a href="http://man.linuxde.net/ntpdate" title="ntpdate">ntpdate</a></li><li><a href="http://man.linuxde.net/dnf" title="dnf">dnf</a></li><li><a href="http://man.linuxde.net/nethogs" title="nethogs">nethogs</a></li><li><a href="http://man.linuxde.net/hping3" title="hping3">hping3</a></li><li><a href="http://man.linuxde.net/trap" title="trap">trap</a></li><li><a href="http://man.linuxde.net/let" title="let">let</a></li><li><a href="http://man.linuxde.net/ifstat" title="ifstat">ifstat</a></li><li><a href="http://man.linuxde.net/blkid" title="blkid">blkid</a></li><li><a href="http://man.linuxde.net/ipcrm" title="ipcrm">ipcrm</a></li><li><a href="http://man.linuxde.net/openssl" title="openssl">openssl</a></li><li><a href="http://man.linuxde.net/chage" title="chage">chage</a></li>            </ul>
        </div>
        <div id="comments">
            <!-- You can start editing here. -->
        </div>
    </div>
    
<div class="sidebar">
    <div class="widget related">
        <div class="hd">
            <h3>相关命令</h3>
        </div>
        <div class="bd">
            <ul>
                            <li><a href="http://man.linuxde.net/at" title="at">at</a></li>
                <li><a href="http://man.linuxde.net/ipcrm" title="ipcrm">ipcrm</a></li>
                <li><a href="http://man.linuxde.net/telint" title="telint">telint</a></li>
                <li><a href="http://man.linuxde.net/w" title="w">w</a></li>
                <li><a href="http://man.linuxde.net/service" title="service">service</a></li>
                <li><a href="http://man.linuxde.net/nohup" title="nohup">nohup</a></li>
                <li><a href="http://man.linuxde.net/pgrep" title="pgrep">pgrep</a></li>
                <li><a href="http://man.linuxde.net/runlevel" title="runlevel">runlevel</a></li>
                <li><a href="http://man.linuxde.net/pmap" title="pmap">pmap</a></li>
                <li><a href="http://man.linuxde.net/pidof" title="pidof">pidof</a></li>
                <li><a href="http://man.linuxde.net/skill" title="skill">skill</a></li>
                <li><a href="http://man.linuxde.net/killall" title="killall">killall</a></li>
                <li><a href="http://man.linuxde.net/batch" title="batch">batch</a></li>
                <li><a href="http://man.linuxde.net/ps" title="ps">ps</a></li>
                <li><a href="http://man.linuxde.net/atq" title="atq">atq</a></li>
                <li><a href="http://man.linuxde.net/init" title="init">init</a></li>
                <li><a href="http://man.linuxde.net/pstree" title="pstree">pstree</a></li>
                <li><a href="http://man.linuxde.net/systemctl" title="systemctl">systemctl</a></li>
                <li><a href="http://man.linuxde.net/atrm" title="atrm">atrm</a></li>
                <li><a href="http://man.linuxde.net/nice" title="nice">nice</a></li>
                </ul>
        </div>
    </div>
    <div id="fixedBox">
        <div id="float" class="dbox">
            <div class="sidebar_search_box">
                <form action="http://man.linuxde.net/">
                    <input type="text" name="s" id="sidebar_search" value="命令搜索" onFocus="if (value =='命令搜索'){value =''}" onBlur="if (value ==''){value='命令搜索'}" />
                </form>
                <form action="http://man.linuxde.net/" method="post">
                    <input type="text" name="ss" id="sidebar_search" value="命令直达" onFocus="if (value =='命令直达'){value =''}" onBlur="if (value ==''){value='命令直达'}" />
                </form>
            </div>
            <div class="widget ad"> 
				<!--
				<script type="text/javascript">
					document.write('<a style="display:none!important" id="tanx-a-mm_22964816_6616438_33736229"></a>');
					tanx_s = document.createElement("script");
					tanx_s.type = "text/javascript";
					tanx_s.charset = "gbk";
					tanx_s.id = "tanx-s-mm_22964816_6616438_33736229";
					tanx_s.async = true;
					tanx_s.src = "http://p.tanx.com/ex?i=mm_22964816_6616438_33736229";
					tanx_h = document.getElementsByTagName("head")[0];
					if(tanx_h)tanx_h.insertBefore(tanx_s,tanx_h.firstChild);
				</script>
				-->
				<a title="linux马克杯" href="http://man.linuxde.net/go.php?http://s.click.taobao.com/t?e=m%3D2%26s%3DbWlkMAz37zAcQipKwQzePDAVflQIoZepLKpWJ%2Bin0XJRAdhuF14FMdl7wKby%2Bc%2FExq3IhSJN6GSAAkUL03DNGUpamuccnq1qrmvEga4qAXJgE5NxN1hjMoxodiO%2Bbyu%2BIYULNg46oBA%3D" target="_blank"><img src="http://112.124.15.130/res/linuxde_man/images/s-ad.jpg" alt="linux马克杯"></a>
            </div>
        </div>
    </div>
</div>
</div>
</div>

<div class="layout layout_ft">
    <div class="footer">
        <div class="copyright">
            <a href="http://man.linuxde.net/" title="Linux命令大全">Linux命令大全</a>
            <em class="sx"></em>
            <a href="http://man.linuxde.net/about_us" title="关于">关于</a>
            <em class="sx"></em>
			收藏本站请使用Ctrl+D
			<em class="sx"></em>
			<a href="http://man.linuxde.net/shell-script" target="_blank">Shell脚本攻略</a>
			<em class="sx"></em>
			<a href="http://man.linuxde.net/docs/shell_regex.html" target="_blank">Shell正则表达式</a>
			<em class="sx"></em>
			共收录528条Linux系统指令
        </div>
        <div class="des">在Linux命令大全（man.linuxde.net）可以查询您所需要的Linux命令教程和相关实例。如果您觉得本站内容对您有所帮助，请推荐给更多需要帮助的人。</div>
        <div class="info">
		赞助商链接：<a href="http://www.1987.name/" target="_blank">组装电脑配置清单2017</a>,　运行在<a href="http://man.linuxde.net/go.php?http://s.click.taobao.com/t?e=m%3D2%26s%3DvTDIoAjlYzMcQipKwQzePCperVdZeJviEViQ0P1Vf2kguMN8XjClAhYdwmcpX19FpcfgoZBkhQKMhxXrDnK9L8ge8W1gUmbKFOH2PYQYfJNjOghCi4TzwOdn1BbglxZYxUhy8exlzcq9AmARIwX9K2Zg%2BdzdQFOwfMRvoxSVDSdLyrb2g0H2G5JcxXijM%2BwneEHpPTskRHknGN7svee8j7S4un6ZpbmpIYULNg46oBA%3D" target="_blank">阿里云服务器</a>！
		</div>
    </div>
</div>
</div>
<!--js-->
<script type="text/javascript" src="http://112.124.15.130/res/linuxde_man/js/sitejs.js"></script>
<script>window._bd_share_config={"common":{"bdSnsKey":{},"bdText":"","bdMini":"2","bdMiniList":false,"bdPic":"","bdStyle":"0","bdSize":"16"},"slide":{"type":"slide","bdImg":"3","bdPos":"right","bdTop":"100"}};with(document)0[(getElementsByTagName('head')[0]||body).appendChild(createElement('script')).src='http://bdimg.share.baidu.com/static/api/js/share.js?v=86835285.js?cdnversion='+~(-new Date()/36e5)];</script>
<div style="display:none"><script src="http://s17.cnzz.com/stat.php?id=5762282&web_id=5762282" language="JavaScript"></script></div>

<!--活动推广
<script type="text/javascript" src="http://112.124.15.130/res/ad/ad.js"></script>
<style type="text/css">
.ad-box{width:240px;height:277px;background:#CCC;position:fixed;right:-240px;bottom:10px;}
.ad-img{position:relative;overflow:hidden;}
.ad-close{width:30px;height:22px;line-height:22px;display:block;top:10px;right:10px;position:absolute;}
</style>
<div class="ad-box">
    <div class="ad-img">
        <a href="http://man.linuxde.net/go.php?http://s.click.taobao.com/t?e=m%3D2%26s%3DuGXlvKu7Iv4cQipKwQzePCperVdZeJviK7Vc7tFgwiFRAdhuF14FMZQUNi%2Bb4SgsMMgx22UI05aAAkUL03DNGT6XZIOJR7dCTAHi9HSRVmdzguIHwwNmQ6UuZxIcp9pfUIgVEmFmgnbDX0%2BHH2IEVaX4VWt66S4EJPwiig1bxLP9BvYCQR6XAr%2BKQ71wHNCAqP8YyUoZZlq4cXg3ii9waXPs9Sj9Qli1np4c65at3FeX3cwyLTlAhj2l4PysJx%2FP" target="_blank"><img src="http://112.124.15.130/res/ad/1111.jpg"></a>
        <a href="javascript:" class="ad-close">关闭</a>
    </div>
</div>
<!--推广结束-->

</body></html>